# services/bookings/service.conf
[program:bookings]
command=python bookings.py
process_name=%(program_name)s
directory=/opt/service
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s/%(program_name)s.log

